# files/etc/init.d/chfs
#!/bin/sh /etc/rc.common
START=99
start() {
    /etc/chfs/chfs --file=/etc/chfs/chfs.ini
}
stop() {
    echo 233 >>/tmp/chfs.log
}
